# deploy/systemd/exzelon-web.service
[Unit]
Description=Exzelon RA Agent — Next.js Frontend
After=network.target

[Service]
Type=simple
User=ra-user
Group=ra-user
WorkingDirectory=/opt/exzelon-ra-agent/frontend
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/node /opt/exzelon-ra-agent/frontend/node_modules/.bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=exzelon-web

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
